# main.cron
0 6 * * * /home/de11an/yasa/project/main.py
